# mcp-backend/backend/supervisord.conf.template
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[inet_http_server]
port=0.0.0.0:${SUPERVISOR_WEB_PORT}
username=${SUPERVISOR_WEB_USERNAME}
password=${SUPERVISOR_WEB_PASSWORD}

[supervisorctl]
serverurl=http://127.0.0.1:${SUPERVISOR_WEB_PORT}

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/supervisord.pid
loglevel=debug
childlogdir=/app/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; MCP server is now managed by the API server directly
; [program:mcp-server]
; command=python /app/mcp_stdio_wrapper.py
; directory=/app
; autostart=true
; autorestart=true
; startsecs=3
; stderr_logfile=/app/logs/mcp-server.err.log
; stdout_logfile=/app/logs/mcp-server.out.log
; stdout_logfile_maxbytes=10MB
; stderr_logfile_maxbytes=10MB
; stdout_logfile_backups=3
; stderr_logfile_backups=3
; environment=MCP_TRANSPORT="stdio",PYTHONUNBUFFERED="1",LOG_LEVEL="DEBUG",PYTHONPATH="/app/ondc-shopping-mcp"
; priority=1

[program:api-server]
command=python /app/api/server.py
directory=/app/api
autostart=true
autorestart=true
startsecs=10
stderr_logfile=/app/logs/api-server.err.log
stdout_logfile=/app/logs/api-server.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",LOG_LEVEL="INFO",MCP_AGENT_CONFIG="/app/mcp_agent.config.yaml"
priority=2

[group:backend]
programs=api-server